# atom_connection.timer
[Unit]
Description=Timer for checking Wire Guard connection

[Timer]
OnUnitInactiveSec=10s
OnBootSec=50s

[Install]
WantedBy=timers.target